# crontabs/cold_backup_crontab
# Определение переменной для каталога проекта
PROJECT_DIR=/home/rus/UTMN/ASOBD/postgres-devops-backup-scripts

# Выполнять скрипт для холодного бэкапа:
# ┌───────────── минута (0 - 59)
# │ ┌───────────── час (0 - 23)
# │ │ ┌───────────── день месяца (1 - 31)
# │ │ │ ┌───────────── месяц (1 - 12)
# │ │ │ │ ┌───────────── день недели (0 - 6) (0=воскресенье)
# │ │ │ │ │
# * * * * *
16 11 * * * cd $PROJECT_DIR && ./scripts/create_cold_backup.sh >> ./logs/create_cold_backup.log 2>&1

